# autostart/rc.local
#!/bin/sh
python3 /home/cantin01/dev/mcbossbot/bossbot.py
